# systemd/ttcoder.socket
[Socket]
ListenStream=/tmp/ttcoder.sock

[Install]
WantedBy=sockets.target
